# sonic_pi_headless/cube-square-music.service
[Unit]
Description=Cube Music Service
After=wireplumber.service

[Service]
WorkingDirectory=/home/pi/sonic-pi-dev
ExecStart=/bin/bash bin/sonic-pi-repl.sh

